# com.yy.roblox-window-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Local alert agent: hourly, pull state/alerts.json and email new GREEN/YELLOW.
     (was 4h; gate D races are hour-scale — the .com pipeline registers and deploys same day)
     Install: fill in GM_SMTP_* below (copy from com.yy.game-monitor.plist), then:
       cp com.yy.roblox-window-monitor.plist ~/Library/LaunchAgents/
       launchctl load ~/Library/LaunchAgents/com.yy.roblox-window-monitor.plist
     Creds live ONLY in this local copy under ~/Library/LaunchAgents (never committed with real values). -->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.yy.roblox-window-monitor</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/sh</string>
    <string>-c</string>
    <string>cd /Users/yy/projects/roblox-window-monitor &amp;&amp; /usr/bin/python3 alert.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/yy/projects/roblox-window-monitor</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>GM_SMTP_USER</key>
    <string>FILL_FROM_game-monitor_plist</string>
    <key>GM_SMTP_PASS</key>
    <string>FILL_FROM_game-monitor_plist</string>
    <key>GM_MAIL_TO</key>
    <string>FILL_FROM_game-monitor_plist</string>
  </dict>

  <key>StartCalendarInterval</key>
  <dict><key>Minute</key><integer>0</integer></dict>

  <key>StandardOutPath</key>
  <string>/Users/yy/Library/Logs/roblox-window-monitor.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/yy/Library/Logs/roblox-window-monitor.log</string>
</dict>
</plist>
