# com.yy.roblox-window-monitor.monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Local monitor agent: runs monitor.mjs every 30 min from the Mac mini (residential IP),
     commits state/ back. Primary since 2026-07-06 — GitHub Actions cron got throttled to ~5h gaps.
     StartCalendarInterval = missed runs coalesce into one catch-up run on wake; RunAtLoad covers reboots.
     No secrets — install verbatim:
       cp com.yy.roblox-window-monitor.monitor.plist ~/Library/LaunchAgents/
       launchctl load ~/Library/LaunchAgents/com.yy.roblox-window-monitor.monitor.plist -->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.yy.roblox-window-monitor.monitor</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/sh</string>
    <string>/Users/yy/projects/roblox-window-monitor/run-monitor.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/yy/projects/roblox-window-monitor</string>

  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Minute</key><integer>15</integer></dict>
    <dict><key>Minute</key><integer>45</integer></dict>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/yy/Library/Logs/roblox-window-monitor.monitor.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/yy/Library/Logs/roblox-window-monitor.monitor.log</string>
</dict>
</plist>
